# aero_pi_ogn_receiver/templates/aero-pi-ogn-rf.service
[Unit]
Description=Aero Pi OGN RF receiver
After=network-online.target time-sync.target
Wants=network-online.target
PartOf=aero-pi-ogn-receiver.target

[Service]
Type=simple
WorkingDirectory=/opt/aero-pi-ogn-receiver/ogn/current
ExecStart=/usr/bin/procServ -f -q -L - --restrict -p /run/aero-pi-ogn-rf.procserv.pid -k ^X --killsig 15 -x ^C -i ^D -c /opt/aero-pi-ogn-receiver/ogn/current 127.0.0.1:50000 ./ogn-rf /etc/aero-pi-ogn-receiver/rtlsdr-ogn.conf
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/aero-pi-ogn-receiver /var/log/aero-pi-ogn-receiver
NoNewPrivileges=true

[Install]
WantedBy=aero-pi-ogn-receiver.target
